# config-files/migration-tasks.conf
[group:migration-tasks]
programs=crond
priority=25

[program:my-program-name]
command = /usr/local/bin/my-script.sh
startsecs = 0
autorestart = false
startretries = 1
